# deploy/backup-bot.service
[Unit]
Description=Telegram DB Backup Bot (Backup Glass)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
EnvironmentFile=-__APP_DIR__/.env
ExecStart=__APP_DIR__/.venv/bin/python __APP_DIR__/main.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=__APP_DIR__/data
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backup-bot

[Install]
WantedBy=multi-user.target
